# build-MX-16.1-SysVinit/Template/COMMON/delete-files.list
root/.bash_history
var/lib/apt/lists/*
var/cache/apt/*.bin
var/cache/apt-show-versions/* 
var/cache/debconf/*-old
var/cache/man/* 
etc/group-
etc/gshadow-
etc/passwd-
etc/shadow-
etc/subgid-
etc/subuid-
etc/machine-id
etc/*.dpkg-dist
etc/*.dpkg-old
etc/issue
etc/issue.net
etc/init.d/*.dpkg-dist
etc/apt/sources.list.d/*.dpkg-dist
etc/apt/sources.list.d/*.dpkg-old
etc/xdg/menus/debian-menu.menu
etc/skel/.xinitrc
var/lib/dpkg/*-old 
var/lib/dbus/machine-id
var/lib/dhcp/*
var/log/*.log
var/log/apt/*.log
var/log/*log
var/log/*mp
var/log/dmesg
#var/cache/flashplugin-nonfree/install*
#var/cache/flashplugin-nonfree/flashplayer*
var/cache/pepperflashplugin-nonfree/google-chrome-*
#usr/share/icons/*/icon-theme.cache 
tmp/*
boot/*.old-dkms
etc/sudoers.d/antixers
etc/rc3.d/S03lightdm
